# laravel-worker.conf
[program:laravel-worker]
command=/usr/bin/php /data/vhosts/usingnet/artisan queue:work redis --sleep=3 --tries=3 --daemon
autostart=true
autorestart=true